# conf/sv-elasticsearch.conf
[program:elasticsearch]
command=/usr/sbin/service elasticsearch start
numprocs=1
autostart=true
autorestart=true
priority=20
