# conf/sv-kibana.conf
[program:kibana]
command=/usr/sbin/service kibana start
numprocs=1
autostart=true
autorestart=true
priority=30
